# systemd/healthcare.service
[Unit]
Description=Healthcare dashboard (docker compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=tsuyoshi
Group=tsuyoshi
WorkingDirectory=/home/tsuyoshi/ghq/github.com/nagamine-git/healthcare
ExecStart=/usr/bin/docker compose --env-file .env.runtime up -d
ExecStop=/usr/bin/docker compose --env-file .env.runtime down
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
